# supervisor-app.ini
[program:app-uwsgi]
command = /usr/local/bin/uwsgi
  --socket /var/run/django.sock
  --chmod-socket=666
  --wsgi-file /code/irdn/wsgi.py
  --static-map /static=/code/static/
  --logto /dev/null
  --processes 4
  --threads 4
  --thunder-lock
  --limit-post 0
  --master-fifo /tmp/uwsgififo
  --buffer-size 32768
  --py-autoreload 3
stopsignal=QUIT

[program:nginx-app]
command = /usr/sbin/nginx

[program:webpack]
command = npm run start
directory = /code/frontend